# deploy/systemd/rob-webhook-dev.service
[Unit]
Description=Rob Dev Webhook Service
After=network.target

[Service]
Type=simple
User=rob
Group=rob
WorkingDirectory=/opt/rob-webhook/app
EnvironmentFile=/opt/rob-webhook/app/.env
Environment=PYTHONPATH=/opt/rob-webhook/app
ExecStart=/opt/rob-webhook/app/.venv/bin/python -m apps.webhook.main
Restart=always
RestartSec=5

# --- sandboxing / hardening ---
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
LockPersonality=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
