# deploy/systemd/rob-webhook.service
[Unit]
Description=Rob Webhook Service
After=network.target

[Service]
Type=simple
User=rob
Group=rob
WorkingDirectory=/opt/rob-webhook/app
EnvironmentFile=/opt/rob-webhook/app/.env
Environment=PYTHONPATH=/opt/rob-webhook/app
ExecStart=/opt/rob-webhook/app/.venv/bin/python -m apps.webhook.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
